# config/systemd/pomo-tui.service
[Unit]
Description=pomo-tui ADHD-focused Pomodoro timer
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%i/bin/pomo-tui daemon
Restart=on-failure
RestartSec=5

# User session service
User=%i
Group=%i

# Environment
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%i
Environment=POMO_TUI_DAEMON_MODE=1

# Working directory
WorkingDirectory=%h

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/pomo-tui %h/.config/pomo-tui

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pomo-tui

[Install]
WantedBy=default.target